# jobs/config/crontab
PATH=/sbin:/usr/sbin:/bin:/usr/bin:/usr/local/bin
* * * * * cd /var/www/html && php artisan schedule:run >> log.txt
